# edge/config/includes.chroot_before_packages/etc/rc.local
#!/bin/bash
rm -f /etc/apt/sources.list
rm -f /etc/rc.local
